# code/runenv-config/fserver/samba/etc-supervisor-conf.d-samba.conf
# ref: https://github.com/silvershell/docker-samba/blob/master/supervisord.conf

[program:smbd]
command=/usr/sbin/smbd --daemon --foreground --log-stdout
redirect_stderr=true

