# code/runenv-config/ms/core/php7/etc-supervisor-conf.d-php.conf
[program:php]
command = /usr/sbin/php-fpm7.0
user = root
autostart = true