# code/runenv-config/fserver/nfs/etc-supervisor-conf.d-nfs.conf
# ref: https://github.com/silvershell/docker-samba/blob/master/supervisord.conf

[program:nfs]
command=nfs-kernel-server
